# etc/systemd/system/wm_timeupdate.timer
[Unit]
Description=WM time update timer
Requires=wm_timeupdate.service

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
